# timers/timertest.service
[Unit]
Description=TimerTest

[Service]
Type=simple
ExecStart=/usr/bin/echo "hello timer-test"
