# timers/timertest.timer
[Unit]
Description=TimerTest

[Timer]
OnActiveSec=2s
