# issue-3166/testsuite.service
[Unit]
Description=Testsuite service
After=multi-user.target

[Service]
ExecStart=/usr/bin/test-fail-on-restart.sh
Type=oneshot
