# issue-2467/test.service
[Unit]
Requires=test.socket
ConditionPathExistsGlob=/var/run/nonexistent

[Service]
ExecStart=/bin/true
