# nspawn-smoke/testsuite.service
[Unit]
Description=Testsuite service
After=multi-user.target

[Service]
ExecStart=/usr/bin/test-nspawn.sh
Type=oneshot
