# unit/assert-virtualization-kvm.service
[Unit]
Description=Service Test AssertVirtualization kvm
AssertVirtualization=kvm

[Service]
ExecStart=/bin/true
Type=oneshot
RemainAfterExit=yes
