# unit/assert_file_not_empty.service
[Unit]
Description=Service Tet AssertFileNotEmpty
AssertFileNotEmpty=/run/assert-file-is-not-empty

[Service]
ExecStart=/bin/true
Type=oneshot
RemainAfterExit=yes
