# unit/condition-virtualization-qemu.service
[Unit]
Description=Service Test ConditionVirtualization qemu
ConditionVirtualization=qemu

[Service]
ExecStart=/bin/true
Type=oneshot
RemainAfterExit=yes
