# unit/condition-file-is-executable.service
[Unit]
Description=Service Test ConditionFileIsExecutable
ConditionFileIsExecutable=/run/file-is-executable

[Service]
ExecStart=/bin/true
Type=oneshot
RemainAfterExit=yes
