# unit/condition-virtualization-lxc.service
[Unit]
Description=Service Test ConditionVirtualization lxc
ConditionVirtualization=lxc

[Service]
ExecStart=/bin/true
Type=oneshot
RemainAfterExit=yes
