# unit/on-failure.service
[Unit]
Description=Service Test OnFailure
OnFailure=on-failure-run.service

[Service]
ExecStart=/bin/false
Type=oneshot
