# unit/on-failure-run.service
[Unit]
Description=Service Test OnFailure

[Service]
ExecStart=/bin/sh -x -c 'touch /run/on-failure-test'
Type=oneshot
